# IoT_SmartHub_Firmware/files_/etc/hotplug.d/usb/1-add
#!/bin/sh
if [ "$ACTION" == "add" ]; then  
/bin/prx800_update.sh &
fi 


