# IoT_SmartHub_Firmware/files_/etc/hotplug.d/iface/99-announce
#!/bin/sh

[ "$ACTION" = ifup ] && {

/etc/init.d/announce stop
/etc/init.d/announce start

}
